# pi/subcensuspi/systemd/subcensuspi.service
[Unit]
# SubCensusPi — the single web-controlled service (Pi §3, §7, §9).
# The dashboard OWNS the one RTL-SDR dongle and switches it between mutually-exclusive modes:
#   off | decode (rtl_433 census -> SQLite/MQTT) | spectrum (live rtl_power waterfall).
# There is no separate collector service — decode runs as a managed subprocess inside this one,
# so two processes never fight over the single radio. The last-selected mode is persisted and
# re-applied on boot (SUBCENSUSPI_RADIO_STATE), so a headless Pi resumes decoding automatically.
Description=SubCensusPi (web dashboard + radio owner: decode/spectrum on one dongle)
Documentation=https://github.com/JamesDavid/SubCensus
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust these lines to your install (defaults match pi/install.sh: $HOME/SubCensus).
User=pi
WorkingDirectory=/home/pi/SubCensus/pi
# Catalog + per-place artifacts the dashboard reads and the decode subprocess writes.
Environment=SUBCENSUSPI_DB=/var/lib/subcensuspi/census.db
Environment=SUBCENSUSPI_PLACE=home
Environment=SUBCENSUSPI_PLACES_DIR=/var/lib/subcensuspi/places
# Collector config for decode mode, and the file that remembers the selected radio mode.
Environment=SUBCENSUSPI_CONFIG=/var/lib/subcensuspi/config.yaml
Environment=SUBCENSUSPI_RADIO_STATE=/var/lib/subcensuspi/radio_state.json
ExecStart=/home/pi/SubCensus/pi/.venv/bin/uvicorn subcensuspi.web.app:app --host 0.0.0.0 --port 8080
# Restart=always is the process safety net. Per-dongle rtl_433 relaunch (Pi §9) is handled INSIDE
# the decode subprocess (supervise_stream, backoff); this restarts the whole owner if it dies.
Restart=always
RestartSec=3
# The decode subprocess needs USB (SDR) access; leave DB/IQ dirs writable to the service user.
StateDirectory=subcensuspi

[Install]
WantedBy=multi-user.target
