# pi/subcensuspi/systemd/subcensuspi-web.service
[Unit]
# SubCensusPi dashboard (Pi §9): read-only FastAPI/uvicorn over the same SQLite catalog.
# Ordered after the collector so the DB exists before the web app opens it.
Description=SubCensusPi web dashboard (FastAPI/uvicorn)
Documentation=https://github.com/JamesDavid/SubCensus
After=network-online.target subcensuspi-collector.service
Wants=network-online.target

[Service]
Type=simple
# Adjust these lines to your install (defaults match pi/install.sh: $HOME/SubCensus).
User=pi
WorkingDirectory=/home/pi/SubCensus/pi
# Point the dashboard at the catalog + per-place artifact dir the collector writes.
Environment=SUBCENSUSPI_DB=/var/lib/subcensuspi/census.db
Environment=SUBCENSUSPI_PLACE=home
Environment=SUBCENSUSPI_PLACES_DIR=/var/lib/subcensuspi/places
ExecStart=/home/pi/SubCensus/pi/.venv/bin/uvicorn subcensuspi.web.app:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
